# dot_config/sv/keymapper/run
#!/bin/sh
exec keymapperd
exec keymapper
