# dot_config/autostart/keymapper.desktop
[Desktop Entry]
Name=keymapper
Comment=Start keymapper
# Icon=keymapper
Exec=keymapper
Terminal=false
Type=Application
NoDisplay=true
