# dot_config/autostart/pipewire.desktop
[Desktop Entry]
Name=PipeWire
Comment=Start PipeWire
Icon=pipewire
Exec=pipewire -c filter-chain.conf
Terminal=false
Type=Application
NoDisplay=true
